# docker/datawake.18-MAY-2016/_etc_supervisor.conf
[supervisord]
nodaemon=true

[program:mongo]
command=mongod

[program:sl-pm]
command=sl-pm 

;[include]
;files = /etc/supervisor/conf.d/*.conf
